# nhsd-backend/jobrunner.service.d/group.conf
[Service]
# nhsd doesn't use our groups
Group=jobrunner
