# release-hatch/release-hatch.timer
[Unit]
Description=Run release-hatch update regularly
Requires=release-release-hatch.service

[Timer]
OnUnitActiveSec=3min
Unit=release-release-hatch.service

[Install]
WantedBy=timers.target
